# config/autostart/gxkb.desktop
[Desktop Entry]
Exec=/usr/bin/gxkb
Name=gxkb
Comment=X11 keyboard indicator and switcher
Type=Application
