# config/autostart/pulseaudio.desktop
[Desktop Entry]
Exec=/usr/bin/pulseaudio --start
Name=pulseaudio
Comment=The PulseAudio Sound System
Type=Application
